# roles/vllm/templates/vllm.service.j2
[Unit]
Description=vLLM Inference Server
After=network.target

[Service]
Type=simple
User={{ vllm_user }}
Group={{ vllm_group }}
SupplementaryGroups=video render
Environment="PATH={{ vllm_install_dir }}/.venv/bin:/usr/local/bin:/usr/bin:/bin"
Environment="HOME={{ vllm_install_dir }}"
ExecStart={{ vllm_install_dir }}/.venv/bin/vllm serve {{ vllm_model }} \
  --host {{ vllm_host }} \
  --port {{ vllm_port }} \
  --gpu-memory-utilization {{ vllm_gpu_memory_utilization }} \
  --tool-call-parser {{ vllm_tool_call_parser }} \
  --enable-auto-tool-choice \
  --attention-backend {{ vllm_attention_backend }} \
  --enable-prefix-caching{% if vllm_api_key is defined %} \
  --api-key {{ vllm_api_key }}{% endif %}

Restart=on-failure
RestartSec=10
TimeoutStartSec=0

[Install]
WantedBy=multi-user.target
